# dashboard/ai-team-dashboard.service
[Unit]
Description=AI Team Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ai-team-dashboard/dashboard
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /root/.openclaw/workspace/ai-team-dashboard/dashboard/server.js
Restart=always
RestartSec=3
StandardOutput=append:/root/.openclaw/workspace/ai-team-dashboard/dashboard/logs/systemd.out.log
StandardError=append:/root/.openclaw/workspace/ai-team-dashboard/dashboard/logs/systemd.err.log

[Install]
WantedBy=multi-user.target
